# crontab
# m h  dom mon dow   command
@reboot (cd /home/selenium/selbot && ./selbot.sh start)
*/10 * * * * (cd /home/selenium/selbot && ./update-index.sh)
